# floorplan-api.service
[Unit]
Description=Floorplan Recognition API Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/path/to/your/app
Environment="PATH=/path/to/your/app/venv/bin"
ExecStart=/path/to/your/app/venv/bin/python app.py
Restart=always

[Install]
WantedBy=multi-user.target 